# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
user=root

# ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log

# ── PHP-FPM ───────────────────────────────────────────────────────────────────
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log

# ── Queue Worker: agent-tasks (high concurrency) ──────────────────────────────
[program:worker-agents]
command=php /var/www/html/artisan queue:work redis --queue=critical,agents --sleep=1 --tries=3 --backoff=30 --max-time=3600
autostart=true
autorestart=true
numprocs=4
process_name=%(program_name)s_%(process_num)02d
priority=20
stdout_logfile=/var/log/supervisor/worker-agents.log
stderr_logfile=/var/log/supervisor/worker-agents-error.log
stopasgroup=true
killasgroup=true

# ── Queue Worker: governance (audit + approvals) ──────────────────────────────
[program:worker-governance]
command=php /var/www/html/artisan queue:work redis --queue=security,governance --sleep=2 --tries=5 --backoff=60 --max-time=3600
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
priority=20
stdout_logfile=/var/log/supervisor/worker-governance.log
stderr_logfile=/var/log/supervisor/worker-governance-error.log
stopasgroup=true
killasgroup=true

# ── Queue Worker: notifications + billing ────────────────────────────────────
[program:worker-default]
command=php /var/www/html/artisan queue:work redis --queue=notifications,billing,workflows,default --sleep=3 --tries=3 --backoff=30 --max-time=3600
autostart=true
autorestart=true
numprocs=2
process_name=%(program_name)s_%(process_num)02d
priority=30
stdout_logfile=/var/log/supervisor/worker-default.log
stderr_logfile=/var/log/supervisor/worker-default-error.log
stopasgroup=true
killasgroup=true
